# systemd/user/plan9.target
[Unit]
Description=Plan9 services
Requires=
After=
Wants=fontsrv.service plumber.service factotum.service

[Install]
WantedBy=default.target
